# Desktop/bot-app/launchd/com.telegram.nsfwbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.telegram.nsfwbot</string>
  <key>WorkingDirectory</key>
  <string>/Users/nishkarshkr/Desktop/myname</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>/Users/nishkarshkr/Desktop/myname/scripts/run_bot.sh</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/nishkarshkr/Library/Logs/telegram-nsfwbot.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nishkarshkr/Library/Logs/telegram-nsfwbot.err.log</string>
</dict>
</plist>
